# download_backup/systemctl/nfs_download.timer
[Unit]
Description=Backup NFS Timer

[Timer]
Unit=nfs_download.service
OnCalendar=23:30:00

[Install]
WantedBy=timers.target